# debian/adaptation-starqlte-configs.bixby-state.service
[Unit]
Description=bixby button gesture service
After=sm.puri.Phosh.service

[Service]
Type=simple
ExecStart=/usr/bin/bixby-state
Restart=always
