# usr/lib/systemd/system/ModemManager.service.d/20-wait.conf
[Service]
ExecStartPre=/usr/bin/sleep 10
